# apps/gadgetron/webapp/gadgetron_web.conf.in
description     "Foobar management daemon"
author          "Alex Smith"

start on started network
start on started startup
stop on stopping network
stop on starting shutdown

console output
kill signal INT

exec su -c "python @CMAKE_INSTALL_PREFIX@/bin/gadgetron_web_app.py @CMAKE_INSTALL_PREFIX@/config/gadgetron_web_app.cfg" gadgetron

respawn

respawn
